# recipes-influx/influx-files/influx-files/etc/systemd/system/autostart.service
[Unit]
Description=Launch my application
After=basic.target

[Service]
Type=simple
ExecStart=/opt/influx/autostart.sh

[Install]
WantedBy=multi-user.target

