# recipes-influx/influx-files/influx-files/etc/systemd/system/lte-ppp.service
[Unit]
Description=LTE PPP connection via Quectel module
After=network.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /var/run/pppd/lock
ExecStart=/usr/sbin/pppd call quectel-ppp nodetach
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
