# recipes-influx/influx-files/influx-files/etc/systemd/system/net-failover.timer
[Unit]
Description=Run network failover script every 30 seconds

[Timer]
OnBootSec=20
OnUnitActiveSec=20

[Install]
WantedBy=timers.target
