# recipes-influx/influx-files/influx-files/etc/systemd/system/wifi_monitor.timer
[Unit]
Description=Run WiFi Monitor periodically

[Timer]
OnBootSec=15s
OnUnitActiveSec=60s

[Install]
WantedBy=timers.target
