# recipes-influx/influx-files/influx-files/etc/profile.d/login.sh
#!/bin/sh
# Put here commands, executed by login 

EDITOR=mcedit   # crontab uses this editor

# check for new firmware version
#/opt/influx/check_firmware_version.sh 1

