# recipes-influx/influx-files/influx-files/etc/profile.d/wpa_supplicant_check.sh
#!/bin/sh

if [ -f /opt/influx/wpa_supplicant.conf.cust ]; then
    mv /opt/influx/wpa_supplicant.conf.cust /etc/wpa_supplicant.conf
fi
